# terraform/setup/upstart/conf.d/airflow_flower.conf
[program:airflow_flower]
command = /usr/local/bin/airflow_flower.bash
user = ubuntu
autostart=true
autorestart=true
