# terraform/setup/airflow-service.cfg
[program:airflow]
command = /home/ubuntu/airflow/airflow_flower.bash

[program:airflow]
process_name=airflow_flower
command = /home/ubuntu/airflow/env/bin/airflow flower
environment=PATH="/home/ubuntu/airflow/env/bin:%(ENV_PATH)s"

[program:airflow]
process_name=airflow_scheduler
command = /home/ubuntu/airflow/env/bin/airflow scheduler
environment=PATH="/home/ubuntu/airflow/env/bin:%(ENV_PATH)s"

[program:airflow]
process_name=airflow_worker
command = /home/ubuntu/airflow/env/bin/airflow worker
environment=PATH="/home/ubuntu/airflow/env/bin:%(ENV_PATH)s"
